# addos/addos-start.service
[Unit]
Description=addos-start

[Service]
ExecStart=/addosinstall.sh

[Install]
WantedBy=multi-user.target
